# configuration/examples/systemd/mailguardian.socket
[Unit]
Description=MailGuardian WSGI Socket

[Socket]
ListenStream=/run/mailguardian.socket

[Install]
WantedBy=sockets.target